# docker/php/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user = root

[program:consumer_bus_search_cache]
command=/bin/ash -c 'sleep 3 && /var/www/bus/bin/console bus-search-cache:consumer --prefetch=1'
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
user=www-data
stdout_logfile_maxbytes = 1GB
stdout_logfile_backups = 3
stderr_logfile_maxbytes = 1GB
stderr_logfile_backups = 3
stderr_logfile = /var/log/supervisor/bus_search-cache_error.log
stdout_logfile = /var/log/supervisor/bus_search-cache_info.log

[program:consumer_bus_payment]
command=/bin/ash -c 'sleep 3 && /var/www/bus/bin/console bus-payment.consumer --prefetch=1'
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
user=www-data
stdout_logfile_maxbytes = 1GB
stdout_logfile_backups = 3
stderr_logfile_maxbytes = 1GB
stderr_logfile_backups = 3
stderr_logfile = /var/log/supervisor/bus_payment_error.log
stdout_logfile = /var/log/supervisor/bus_payment_info.log

[program:consumer_bus_periodically_update]
command=/bin/ash -c 'sleep 3 && /var/www/bus/bin/console bus-periodically-update.consumer --prefetch=1'
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
user=www-data
stdout_logfile_maxbytes = 1GB
stdout_logfile_backups = 3
stderr_logfile_maxbytes = 1GB
stderr_logfile_backups = 3
stderr_logfile = /var/log/supervisor/bus_periodically_update_error.log
stdout_logfile = /var/log/supervisor/bus_periodically_update_info.log

[program:consumer_bus_periodically_update]
command=/bin/ash -c 'sleep 3 && /var/www/bus/bin/console bus-periodically-update.consumer --prefetch=1'
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
user=www-data
stdout_logfile_maxbytes = 1GB
stdout_logfile_backups = 3
stderr_logfile_maxbytes = 1GB
stderr_logfile_backups = 3
stderr_logfile = /var/log/supervisor/bus_periodically_update_error.log
stdout_logfile = /var/log/supervisor/bus_periodically_update_info.log

[program:consumer_bus_elastic_update_sh]
command=/bin/ash -c 'sleep 3 && /var/www/bus/bin/console bus-elastic-update.consumer --prefetch=1 --datacenter=sh'
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
user=www-data
stdout_logfile_maxbytes = 1GB
stdout_logfile_backups = 3
stderr_logfile_maxbytes = 1GB
stderr_logfile_backups = 3
stderr_logfile = /var/log/supervisor/bus_elastic_update_sh_error.log
stdout_logfile = /var/log/supervisor/bus_elastic_update_sh_info.log

[program:consumer_bus_elastic_update_rd]
command=/bin/ash -c 'sleep 3 && /var/www/bus/bin/console bus-elastic-update.consumer --prefetch=1 --datacenter=rd'
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
user=www-data
stdout_logfile_maxbytes = 1GB
stdout_logfile_backups = 3
stderr_logfile_maxbytes = 1GB
stderr_logfile_backups = 3
stderr_logfile = /var/log/supervisor/bus_elastic_update_rd_error.log
stdout_logfile = /var/log/supervisor/bus_elastic_update_rd_info.log

[program:consumer_bus_async_high]
command=/bin/ash -c 'sleep 3 && /var/www/bus/bin/console bus-async-event.consumer --priority=high'
process_name=%(program_name)s_%(process_num)02ds
numprocs=1
autostart=true
autorestart=true
user=www-data
stdout_logfile_maxbytes = 1GB
stdout_logfile_backups = 3
stderr_logfile_maxbytes = 1GB
stderr_logfile_backups = 3
stderr_logfile = /var/log/supervisor/bus_async_high_error.log
stdout_logfile = /var/log/supervisor/bus_async_high_info.log

[program:consumer_bus_async_medium]
command=/bin/ash -c 'sleep 3 && /var/www/bus/bin/console bus-async-event.consumer --priority=medium'
process_name=%(program_name)s_%(process_num)02ds
numprocs=1
autostart=true
autorestart=true
user=www-data
stdout_logfile_maxbytes = 1GB
stdout_logfile_backups = 3
stderr_logfile_maxbytes = 1GB
stderr_logfile_backups = 3
stderr_logfile = /var/log/supervisor/bus_async_medium_error.log
stdout_logfile = /var/log/supervisor/bus_async_medium_error.log

[program:consumer_bus_async_low]
command=/bin/ash -c 'sleep 3 && /var/www/bus/bin/console bus-async-event.consumer --priority=low'
process_name=%(program_name)s_%(process_num)02ds
numprocs=1
autostart=true
autorestart=true
user=www-data
stdout_logfile_maxbytes = 1GB
stdout_logfile_backups = 3
stderr_logfile_maxbytes = 1GB
stderr_logfile_backups = 3
stderr_logfile = /var/log/supervisor/bus_async_low_error.log
stdout_logfile = /var/log/supervisor/bus_async_low_error.log
